# systemd/secops-buddy.timer
[Unit]
Description=SecOps Buddy timer (deprecated)

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
