# systemd/secops-buddy.service
[Unit]
Description=SecOps Buddy (agent + bot)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/secops-buddy
Environment=SECOPS_BUDDY_ROOT=/opt/secops-buddy
EnvironmentFile=/etc/secops-buddy/.env
ExecStart=/opt/secops-buddy/.venv/bin/python /opt/secops-buddy/run.py --config /etc/secops-buddy/config.yml
User=root
Restart=always
RestartSec=2
KillSignal=SIGTERM
TimeoutStopSec=25

[Install]
WantedBy=multi-user.target
